# systemd/.config/systemd/user/unclutter.service
[Unit]
Description=Unclutter
[Service]
ExecStart=/usr/bin/unclutter -idle 2 -noevents
[Install]
WantedBy=default.target
