# systemd/.config/systemd/user/xscreensaver.service
[Unit]
Description=XScreenSaver
[Service]
ExecStart=/usr/bin/xscreensaver -no-capture-stderr
[Install]
WantedBy=default.target
